# winwing-dcs-bridge.service
[Unit]
Description=WinWing DCS Telemetry Bridge
After=network.target

[Service]
Type=simple
# Your path here
ExecStart=/usr/bin/python3 -u %h[path_to_script]
Restart=always
RestartSec=5

# The bridge idles at 1Hz when DCS isn't running, so it's safe to keep always-on
# Logs go to journalctl: journalctl --user -u winwing-dcs-bridge -f

[Install]
WantedBy=default.target
